# news-monitor.service
[Unit]
Description=News Monitor Application
Documentation=https://github.com/your-repo/news-monitor
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/news-monitor

# 环境变量
Environment="PATH=/path/to/news-monitor/venv/bin"
EnvironmentFile=/path/to/news-monitor/.env

# 启动命令
ExecStart=/path/to/news-monitor/venv/bin/python monitor_app.py

# 重启策略
Restart=always
RestartSec=10

# 资源限制
LimitNOFILE=65535
Nice=5

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=news-monitor

[Install]
WantedBy=multi-user.target
